# fonera/etc/crontabs/root
*  */12 * * * /www/cgi-bin/set_time.sh
*/30  * * * * /etc/init.d/ntpclient
